# timer/pass-close@.service
[Unit]
Description=pass-close %i

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/pass close --verbose %i
